# odysseus-bot.service
[Unit]
Description=Odysseus Telegram Bot
Documentation=https://github.com/GODENERGON/Apextrading
After=network-online.target docker.service
Wants=network-online.target
ConditionHost=godener-NUC8i7BEH

[Service]
Type=simple
User=godener
Group=godener
WorkingDirectory=/home/godener/odysseus/bot

# Environment from .env
EnvironmentFile=/home/godener/odysseus/bot/.env

# Python and dependencies
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /home/godener/odysseus/bot/main.py

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=odysseus-bot

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
